# thehive-qradar.timer
[Unit]
Description=Ticket  Append Service Every 3 Hours


[Timer]
OnCalendar=*-*-* 00/1:00:00
Persistent=true

[Install]
WantedBy=timers.target
